# examples/server-nginx/.docker/services/nextjs/finish
#!/usr/bin/execlineb -S1
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services


# stop docker container if this service crashes
