# examples/server-nginx/.docker/services/nextjs/run
#!/usr/bin/execlineb -P

with-contenv
s6-envdir env
s6-setuidgid app
cd /usr/local/app
yarn start
