# examples/server-nginx/.docker/services/nginx/run
#!/usr/bin/execlineb -P

with-contenv
s6-envdir env
/usr/sbin/nginx -g "daemon off;"
